# cd-image/source/root/etc/profile.d/gog.sh
export vblank_mode=0
export PATH=$PATH:/usr/lib64/mesa/bin/xdemos
export PATH=$PATH:/usr/lib64/mesa/bin/egl/opengl
export PATH=$PATH:/usr/sbin
export EDITOR=mcedit